# quadlet/calculator3d-backend.container
[Unit]
Description=Calculator3D Backend (FastAPI)
After=network-online.target

[Container]
ContainerName=calculator3d-backend
Image=localhost/calculator3d-backend:latest
Network=calculator3d.network
Volume=calculator3d-data.volume:/app/data:Z
EnvironmentFile=__DEPLOY_PATH__/.env
Environment=DATABASE_URL=sqlite+aiosqlite:///./data/calculator3d.db
Environment=ALGORITHM=HS256
Environment=ACCESS_TOKEN_EXPIRE_MINUTES=1440

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
