# quadlet/calculator3d-frontend.container
[Unit]
Description=Calculator3D Frontend (React + Nginx)
After=calculator3d-backend.service

[Container]
ContainerName=calculator3d-frontend
Image=localhost/calculator3d-frontend:latest
PublishPort=3000:80
Network=calculator3d.network

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
