# quadlet/cfs-frontend.container
[Unit]
Description=collectors-forge-studio Frontend (React + Nginx)
After=cfs-backend.service

[Container]
ContainerName=cfs-frontend
Image=localhost/cfs-frontend:latest
PublishPort=3000:80
Network=cfs.network

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
