# quadlet/cfs-tracker.container
[Unit]
Description=Collector's Forge Tracker - Seguimiento de pedidos con Playwright
After=network-online.target cfs-postgres.service
Requires=cfs-postgres.service

[Container]
ContainerName=cfs-tracker
Image=localhost/cfs-tracker:latest
Network=cfs.network:alias=tracker
EnvironmentFile=__ENV_FILE__
HealthCmd=curl -sf http://localhost:8002/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
