# quadlet/cfs-minio.container
# NOTA: Este Quadlet va en turtleStorage, NO en turtleServer.
# Instalar manualmente en /home/turtlestorage/.config/containers/systemd/
# Luego: systemctl --user daemon-reload && systemctl --user start cfs-minio

[Unit]
Description=Collector's Forge MinIO - Object storage para el Vault
After=network-online.target

[Container]
ContainerName=cfs-minio
Image=quay.io/minio/minio:latest
PublishPort=9000:9000
PublishPort=9001:9001
Volume=cfs-minio-data.volume:/data
EnvironmentFile=/home/turtlestorage/.env
Exec=server /data --console-address ":9001"
HealthCmd=curl -sf http://localhost:9000/minio/health/live
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
